# scripts/django.conf
[program:celery-worker]
process_name=%(program_name)s_%(process_num)02d
command=python3 manage.py runserver 0.0.0.0:887
directory=/home/python/cloudpanel/
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/python/cloudpanel/django.log